# debian/postinst
#!/bin/sh

update-rc.d wb-rules defaults
invoke-rc.d wb-rules start
